# lib/generators/script_generator/templates/startup_server.erb
#! /bin/bash
### BEGIN INIT INFO
# Provides:          <%= Rails.application.engine_name.titleize.chomp(' Application') %>
# Required-Start:    $local_fs $remote_fs $network $syslog redis-server memcached
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start Up application server
# Description:       Start Up application server
### END INIT INFO

NAME=<%= Rails.application.engine_name.chomp('_application') %>
DESC="<%= Rails.application.engine_name.titleize.chomp(' Application') %> application server"

APP_PATH=<%= "#{Setting["deployment"]["path"]}/#{Setting["deployment"]["app_name"]}/current" || Rails.root %>
BUNDLE_CMD=<%= Setting["deployment"]["bundle_wrapper_cmd"] || "bundle" %>

PID=$APP_PATH/tmp/pids/unicorn.pid
RESQUE_PID=$APP_PATH/tmp/pids/resque.pid

DAEMON_OPTS="-c $APP_PATH/config/unicorn.rb -E production -D"
RESQUE_OPTS="QUEUE=* RAILS_ENV=production PIDFILE=$RESQUE_PID"

case "$1" in
  start)
        CD_TO_APP_DIR="cd $APP_PATH"
        START_DAEMON_PROCESS="$BUNDLE_CMD exec unicorn_rails $DAEMON_OPTS"
        # Not using BACKGROUND=yes cause it no log work anymore since resque 1.20
        START_RESQUE_PROCESS="$BUNDLE_CMD exec rake environment resque:work $RESQUE_OPTS"

        echo -n "Starting $DESC: "
        $CD_TO_APP_DIR > /dev/null 2>&1 && $START_DAEMON_PROCESS && $START_RESQUE_PROCESS 2>&1 >> $APP_PATH/log/resque.log &
        echo "$NAME."
        ;;
  stop)
        echo -n "Stopping $DESC: "
        kill -QUIT `cat $PID`
        kill -QUIT `cat $RESQUE_PID`
        echo "$NAME."
        ;;
  restart)
        echo -n "Restarting $DESC: "
        kill -USR2 `cat $PID`
        echo "$NAME."
        ;;
  reload)
        echo -n "Reloading $DESC configuration: "
        kill -HUP `cat $PID`
        echo "$NAME."
        ;;
  *)
        echo "Usage: $NAME {start|stop|restart|reload}" >&2
        exit 1
        ;;
esac

exit 0
